# supervisord.conf
[supervisord] 
nodaemon=false

[program:php-fpm] 
command=/usr/sbin/php-fpm -D 

[program:nginx] 
command=/usr/sbin/nginx -DFOREGROUND -k start
